# etc/s6-overlay/s6-rc.d/sshd/run
#!/usr/bin/execlineb -P
with-contenv
/usr/sbin/sshd -D
