# etc/s6-overlay/s6-rc.d/ngrok/run
#!/command/execlineb -P
with-contenv
/ngrok tcp --remote-addr=${NGROK_DOMAIN} 22
