# etc/s6-overlay/s6-rc.d/ngrok/run
#!/command/with-contenv sh
/ngrok config add-authtoken "$Ngrok"
/ngrok tcp --remote-addr=${NGROK_DOMAIN} 22
