# etc/s6-overlay/s6-rc.d/custom-service/run
#!/command/execlineb -P
if -t { test -f /volume/devbox/run-service.sh } /usr/bin/bash /volume/devbox/run-service.sh
