# etc/s6-overlay/s6-rc.d/custom-service/run
#!/usr/bin/execlineb -P
with-contenv
if -t { s6-test -f /volume/devbox/run-service.sh } /volume/devbox/run-service.sh
