# etc/s6-overlay/s6-rc.d/ngrok-init/up
#!/command/with-contenv sh
/ngrok config add-authtoken ${Ngrok}
